# php.cron
*/1 * * * * curl http://web/cron.php --header 'Host: admin.outline' > /dev/null
